# systemd/docker_swarm_service.service
[Unit]
Description=Docker Swarm Service Manager
After=network.target docker.service
Requires=docker.service

[Service]
Type=simple
ExecStart=docker ps
ExecStop=docker rm -f $(docker ps | grep -i xcat_stack | awk {'print $1})
ExecReload=docker rm -f $(docker ps | grep -i xcat_stack | awk {'print $1})
Restart=always
RestartSec=10
EnvironmentFile=-/etc/sysconfig/docker_swarm_service

[Install]
WantedBy=multi-user.target
